# debian/postinst
#!/bin/sh

set -e

update-rc.d firstboot defaults 99
